# crontab
0 */3 * * * python3 /pec.py >> /logfile
